# systemd/honeypot-ssh.service
[Unit]
Description=Honeypot SSH (paramiko)
After=network.target

[Service]
Type=simple
User=kali
Group=kali
WorkingDirectory=/home/kali/honeypots/ssh
ExecStart=/home/kali/honeypots/venv/bin/python /home/kali/honeypots/ssh/honeypot_ssh.py
Restart=on-failure
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
